# files/Debian/cron.d/local-zabbix-check-upgrades
# Distributed via ansible - mit.zabbix-agent.check-upgrades

# 2020-06-29: Added "Failed to check for VM: Transport endpoint is not connected"
# 2022-11-14: Disabled stdout and stderr
*/5 * * * *     root    /opt/bin/zabbix-check-upgrades >/dev/null 2>&1
